# dist/markco.ini
[supervisord]
nodaemon=true
environmant=HOME=/var/lib/markco

[program:markco-exe]
command=markco-exe
user=markco
autorestart=true

[program:caddy]
command=caddy -conf /etc/Caddyfile
user=markco
autorestart=true
